# discord_bot.service
[Unit]
Description=xkcd discord bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/jblitzar/xkcd-rag
ExecStart=/bin/sh -c 'exec /home/jblitzar/.local/bin/uv run bot.py'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target